# rootfs/etc/s6-overlay/s6-rc.d/wordpress-setup/up
/etc/s6-overlay/scripts/wordpress-setup.sh
