# deploy/wordpress-bedrock.service
[Unit]
Description=WordPress Bedrock Compose Stack with Traefik TLS
Requires=docker.service
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/srv/wordpress-bedrock
EnvironmentFile=-/srv/wordpress-bedrock/.env
ExecStart=/usr/bin/docker compose -f docker-compose.yaml -f docker-compose.tls.yml up -d --build
ExecStop=/usr/bin/docker compose -f docker-compose.yaml -f docker-compose.tls.yml down
ExecReload=/usr/bin/docker compose -f docker-compose.yaml -f docker-compose.tls.yml up -d --build
TimeoutStartSec=0

[Install]
WantedBy=multi-user.target
